# nextcloud/crontab
# Nextcloud
*/5 * * * * docker exec --user www-data nextcloud-app php -f cron.php > /dev/null 2>&1
*/15 * * * * docker exec --user www-data nextcloud-app ./occ preview:pre-generate > /dev/null 2>&1
